# rootfs_files/usr/lib/udev/rules.d/89-mos_automount.rules
# MOS: Check if device needs to be auto mounted
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sd*|nvme*|vd*", RUN+="/usr/local/bin/mos-udev_automount"
